# ops/chloe.service
[Unit]
Description=Chloe 2.0
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/chloe
EnvironmentFile=/opt/chloe/.env
ExecStart=/opt/chloe/.venv/bin/uvicorn chloe.app:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
